# deploy/systemd/budget-app.service
[Unit]
Description=Budget App Next.js service
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/budget-app/app
EnvironmentFile=/etc/budget-app/budget-app.env
ExecStart=/usr/bin/env npm run start
Restart=always
RestartSec=5
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
